# supervisor/templates/process.conf.jinja
; This file is managed by Salt. Do not edit.
; Source: {{ source }}
[program:{{ id }}]
command={{ command }}
user={{ user }}
autostart={{ autostart }}
autorestart={{ autorestart }}
startsecs={{ startsecs }}
startretries={{ startretries }}
exitcodes={{ exitcodes }}
stopsignal={{ stopsignal }}
stopwaitsecs={{ stopwaitsecs }}
redirect_stderr={{ redirect_stderr }}
numprocs={{ numprocs }}
